# ubuntu22/ansible/templates/update-motd.d/99-aracpac
#!/bin/sh
/usr/bin/neofetch --source /etc/aracpac-ascii